# deploy/custom-dispatcher/cimasim-custom-dispatcher.service
[Unit]
Description=CimaSim custom Xyce rootless dispatcher
After=network-online.target

[Service]
Type=simple
User=cimasim-runner
Group=cimasim-runner
SupplementaryGroups=cimasim-custom-spool
Environment=HOME=/var/lib/cimasim-runner
Environment=XDG_RUNTIME_DIR=/run/user/997
EnvironmentFile=/etc/cimasim/custom-dispatcher.env
ExecStart=/opt/cimasim/custom-dispatcher/bin/cimasim-custom-dispatcher --spool-root ${CIMASIM_CUSTOM_SPOOL_ROOT} --runner-image-id ${CIMASIM_CUSTOM_RUNNER_IMAGE_ID} --poll-seconds 2
Restart=on-failure
RestartSec=5s
KillMode=mixed
TimeoutStopSec=20s
UMask=0007
Delegate=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RemoveIPC=yes
ReadWritePaths=/var/lib/cimasim-custom-spool /var/lib/cimasim-runner /run/user/997

[Install]
WantedBy=multi-user.target
